# infra/systemd/moq-cert.service
[Unit]
Description=MoQ Certificate Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/moq-relay/env

# Install certbot and plugin if needed
ExecStartPre=/bin/sh -c 'if ! command -v certbot &> /dev/null; then nix profile install nixpkgs#certbot nixpkgs#certbot-dns-google; fi'

# Request wildcard cert for moq.hang.live and node-specific cert
ExecStart=/bin/sh -c 'certbot certonly \
  --dns-google \
  --dns-google-credentials /etc/moq-relay/gcp-dns.json \
  --non-interactive \
  --agree-tos \
  --email admin@hang.live \
  -d moq.hang.live \
  -d ${NODE_NAME}.moq.hang.live \
  --deploy-hook "systemctl reload-or-restart moq-relay.service"'

# Mark that we've run at least once
ExecStartPost=/usr/bin/touch /var/lib/moq-relay/.cert-initialized

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
