# infra/relay/moq-relay.service
[Unit]
Description=MoQ Relay Server
After=network-online.target moq-cert.service
Wants=network-online.target
Requires=moq-cert.service

[Service]
Type=notify
User=root
WorkingDirectory=/var/lib/moq
Environment="RUST_LOG=trace"
Environment="RUST_BACKTRACE=1"

ExecStart=/var/lib/moq/relay/bin/moq-relay \
  --listen [::]:443 \
  --tls-cert /etc/letsencrypt/live/moq.hang.live/fullchain.pem \
  --tls-key /etc/letsencrypt/live/moq.hang.live/privkey.pem \
  --auth-key /var/lib/moq/root.jwk \
  --auth-public anon \
  --web-https-listen [::]:443 \
  --web-https-cert /etc/letsencrypt/live/moq.hang.live/fullchain.pem \
  --web-https-key /etc/letsencrypt/live/moq.hang.live/privkey.pem \
  --cluster-root use.moq.hang.live \
  --cluster-node %H \
  --cluster-token /var/lib/moq/cluster.jwt

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening - relaxed to allow Nix to work
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
