# infra/relay/vacuum.service
[Unit]
Description=Vacuum systemd journal logs and Nix store
Documentation=man:journalctl(1) man:nix-collect-garbage(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=7d --vacuum-size=500M
ExecStart=/usr/bin/nix-collect-garbage --delete-older-than 7d

StandardOutput=journal
StandardError=journal
